# deploy/millisfera-deploy.service
[Unit]
Description=Sync Millisfera from GitHub to the live Coolify source mount
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/coolify-source/millisfera-az
ExecStart=/usr/bin/env bash /opt/coolify-source/millisfera-az/deploy/pull-deploy.sh
